# supervisor/conf.d/decameron.conf
[program:decameron]
command=docker-compose up -d
directory=/root/docker
redirect_stderr=true
autostart=true
autorestart=true
priority=10